# supervisord.conf
[supervisord]
nodaemon=true

[program:rails]
command=rails server

